# host/systemd/dca.service
[Unit]
Description=DCA

[Service]
Type=simple
ExecStart=/opt/bin/send.sh
