# host/systemd/dca.timer
[Unit]
Description=DCA Info

[Timer]
OnUnitActiveSec=60
OnBootSec=1min
Unit=dca.service

[Install]
WantedBy=multi-user.target
